# Team1/ros2/systemd/v2x-alert-client-docker.service

[Unit]
Description=V2X Alert Client (Docker, ROS2 publisher)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
Restart=always
RestartSec=2

# --- 환경변수 (필요 시 수정) ---
Environment=MCAST=239.20.20.20
Environment=PORT=5520
Environment=ROS_DOMAIN_ID=10
Environment=V2X_KEY=

# 이전 컨테이너 정리
ExecStartPre=-<DOCKER_BIN> rm -f v2x-client

# 컨테이너 실행
ExecStart=<DOCKER_BIN> run --rm --name v2x-client \
  --net=host --ipc=host \
  -e ROS_DOMAIN_ID=${ROS_DOMAIN_ID} \
  -e V2X_KEY=${V2X_KEY} \
  -e MCAST=${MCAST} -e PORT=${PORT} \
  -v <V2X_DIR>:<V2X_DIR> \
  ros:jazzy-ros-base \
  /bin/bash -lc 'source /opt/ros/jazzy/setup.bash; \
    python3 <V2X_DIR>/client.py \
      --mcast "${MCAST}" \
      --port "${PORT}" \
      --log <V2X_DIR>/alerts.csv \
      --ros2 \
      ${V2X_KEY:+--hmac-key "$V2X_KEY"}'

ExecStop=<DOCKER_BIN> stop v2x-client

[Install]
WantedBy=multi-user.target
