# Team1/ros2/systemd/v2x-gateway.service

[Unit]
Description=V2X Gateway (UDP -> ROS2 /v2x/alert in container)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_GROUP>
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=2

# 래퍼 없이 직접 실행
ExecStart=/bin/bash -lc '\
  <DOCKER_BIN> rm -f v2x-gateway >/dev/null 2>&1 || true; \
  exec <DOCKER_BIN> run --name v2x-gateway \
    --net=host --ipc=host \
    -e ROS_DOMAIN_ID=10 \
    -e PYTHONUNBUFFERED=1 \
    -v <PROJECT_ROOT>:/ws \
    -v /home/pi/v2v:<V2X_DIR>:ro \
    -w /ws \
    ros:jazzy-ros-base \
    bash -lc "source /opt/ros/$ROS_DISTRO/setup.bash && source /ws/install/setup.bash || true; python3 -u <V2X_DIR>/client.py --ros2 --topic /v2x/alert --qos reliable --depth 10" \
'

ExecStop=<DOCKER_BIN> stop v2x-gateway
ExecStopPost=<DOCKER_BIN> rm -f v2x-gateway

[Install]
WantedBy=multi-user.target
